# prometheus/rootfs/etc/services.d/prometheus/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: prometheus
# Runs the prometheus
# ==============================================================================

declare ingress_interface
declare ingress_port

bashio::log.info "Starting Prometheus..."

ingress_port=$(bashio::addon.ingress_port)
ingress_interface=$(bashio::addon.ip_address)

# Run the prometheus
exec /bin/prometheus --config.file /etc/prometheus/prometheus.yml --web.external-url http://$(ingress_interface):$(ingress_port)