# prometheus/rootfs/etc/services.d/prometheus/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: prometheus
# Runs the prometheus
# ==============================================================================
declare -a options
declare ingress_interface
declare ingress_port
declare ingress_entry

ingress_port=$(bashio::addon.ingress_port)
ingress_interface=$(bashio::addon.ip_address)

options+=(--config.file="/etc/prometheus/prometheus.yml")
options+=(--storage.tsdb.retention.time="15d")

bashio::log.info "Starting Prometheus..."

# Run the prometheus
exec /bin/prometheus "${options[@]}"