# prometheus/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Prometheus
# Runs the nginx daemon
# ==============================================================================

bashio::log.info "Starting NGinx..."

exec nginx