# prometheus/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Prometheus
# Runs the nginx daemon
# ==============================================================================

# Wait for Prometheus to start
bashio::net.wait_for 9090

bashio::log.info "Starting Nginx..."

exec nginx